# infra/systemd/atelier-manager.service
[Unit]
Description=L'atelier Manager API
Documentation=https://github.com/frak-id/frak-sandbox
After=network.target atelier-network.service
Wants=atelier-network.service
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
User=atelier
Group=atelier
SupplementaryGroups=kvm
WorkingDirectory=/opt/atelier
KillMode=process
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
Environment=SANDBOX_MODE=production
Environment=DATA_DIR=/opt/atelier/data
Environment=MIGRATIONS_DIR=/opt/atelier/drizzle
ExecStart=/usr/local/bin/bun server.js

Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=atelier-manager

[Install]
WantedBy=multi-user.target
